# cert_check/systemd/check-certs.timer
[Unit]
Description=Run TLS certificate checker periodically

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target


